# system/.config/systemd/user/pull_oracle.timer
[Unit]
Description=Pull from oracle daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=pull_oracle.service

[Install]
WantedBy=timers.target
